# units/apiary-tmux.service.template
[Unit]
Description=Apiary — tmux session "apiary" (long-lived, outlives ttyd)

[Service]
Type=oneshot
RemainAfterExit=yes
KillMode=process
Environment=TERM=xterm-256color
ExecStart=/bin/sh -c '${TMUX_BIN} has-session -t apiary 2>/dev/null || ${TMUX_BIN} new-session -d -s apiary'
ExecStop=${TMUX_BIN} kill-session -t apiary

[Install]
WantedBy=default.target
